# docker/cubox-i/build/etc/local.d/50-retrieve-firmware.start
mkdir /lib/firmware/imx/sdma -p
mkdir /lib/firmware/brcm -p
mkdir /lib/firmware/ti-connectivity -p

if [ ! -e "/lib/firmware/vpu_fw_imx6q.bin" ]
then
	wget http://www.nxp.com/lgfiles/NMG/MAD/YOCTO/firmware-imx-8.0.bin
	dd if=firmware-imx-8.0.bin bs=37180 skip=1 | tar xj
	cp firmware-imx-8.0/firmware/vpu/vpu_fw_imx6q.bin /lib/firmware
fi

if [ ! -e "/lib/firmware/imx/sdma/sdma-imx6q.bin" ]
then
	wget https://git.kernel.org/pub/scm/linux/kernel/git/firmware/linux-firmware.git/plain/imx/sdma/sdma-imx6q.bin
	cp sdma-imx6q.bin /lib/firmware/imx/sdma/sdma-imx6q.bin
fi

if [ ! -e "/lib/firmware/brcm/brcmfmac4330-sdio.bin" ]
then
	wget https://git.kernel.org/pub/scm/linux/kernel/git/firmware/linux-firmware.git/plain/brcm/brcmfmac4330-sdio.bin
	cp brcmfmac4330-sdio.bin /lib/firmware/brcm/brcmfmac4330-sdio.bin
fi
if [ ! -e "/lib/firmware/brcm/brcmfmac4330-sdio.txt" ]
then
	wget https://github.com/Freescale/meta-freescale-3rdparty/raw/master/recipes-bsp/broadcom-nvram-config/files/cubox-i/brcmfmac4330-sdio.txt
	cp brcmfmac4330-sdio.txt /lib/firmware/brcm/brcmfmac4330-sdio.txt
fi

if [ ! -e "/lib/firmware/brcm/brcmfmac4329-sdio.bin" ]
then
	wget https://git.kernel.org/pub/scm/linux/kernel/git/firmware/linux-firmware.git/plain/brcm/brcmfmac4329-sdio.bin
	cp brcmfmac4329-sdio.bin /lib/firmware/brcm/brcmfmac4329-sdio.bin
fi
if [ ! -e "/lib/firmware/brcm/brcmfmac4329-sdio.txt" ]
then
	wget https://github.com/Freescale/meta-freescale-3rdparty/raw/master/recipes-bsp/broadcom-nvram-config/files/cubox-i/brcmfmac4329-sdio.txt
	cp brcmfmac4329-sdio.txt /lib/firmware/brcm/brcmfmac4329-sdio.txt
fi

if [ ! -e "/lib/firmware/ti-connectivity/TIInit_11.8.32.bts" ]
then
	wget https://github.com/TI-ECS/bt-firmware/raw/master/TIInit_11.8.32.bts
	cp TIInit_11.8.32.bts /lib/firmware/ti-connectivity/TIInit_11.8.32.bts
fi

if [ ! -e "/lib/firmware/ti-connectivity/wl1271-nvs.bin" ]
then
	wget https://github.com/TI-OpenLink/firmwares/raw/master/ti-connectivity/wl1271-nvs.bin
	cp wl1271-nvs.bin /lib/firmware/ti-connectivity/wl1271-nvs.bin
fi
